# deploy/supervisor-unikosana-scheduler.conf
; Supervisor program — Laravel scheduler for unikosana (chapter).
; Runs `schedule:work`, which invokes the scheduler every minute in the
; foreground (no crontab needed). This is what fires the daily backups:
; backup:clean (01:00) and backup:run (01:30) via spatie/laravel-backup.
;
; Install with root:
;   sudo cp deploy/supervisor-unikosana-scheduler.conf /etc/supervisor/conf.d/unikosana-scheduler.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status unikosana-scheduler
;
; After a deploy:  sudo supervisorctl restart unikosana-scheduler

[program:unikosana-scheduler]
command=/usr/bin/php /home/unikosana/unikosana/artisan schedule:work
directory=/home/unikosana/unikosana
user=deploy
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/unikosana/unikosana/storage/logs/scheduler.log
